# docker/conf/crontab
SHELL=/bin/bash
BASH_ENV=/container.env
*/1 * * * * /bin/bash /root/partyman/bin/cron.sh >> /root/partyman/cron.log 2>&1
# extra line
